# docker/supervisor/supervisor_php_fpm.conf
[program:php-fpm]
service php-fpm7.0 start
